# systemd/radiostream-ffmpeg.service
# /etc/systemd/system/radiostream-ffmpeg.service
# Captures PulseAudio output and streams to Icecast

[Unit]
Description=FFmpeg audio capture for RadioStream
After=radiostream-pulseaudio.service icecast2.service
Requires=radiostream-pulseaudio.service

[Service]
Type=simple
ExecStart=/usr/bin/ffmpeg \
  -fflags nobuffer \
  -flags low_delay \
  -f pulse -i spotify_out.monitor \
  -acodec aac -ab 192k -ar 44100 \
  -flush_packets 1 \
  -frag_duration 100000 \
  -content_type audio/aac \
  -f adts icecast://source:hackme@127.0.0.1:8000/isla
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
